# others/etc/systemd/system/php-fpm.service.d/override.conf
[Service]
LimitNPROC=65535
LimitNOFILE=100000
LimitSTACK=infinity
PrivateTmp=false
